# services.d/health/run
#!/usr/bin/execlineb -P
/healthcheck -name animated-octopus
